# src/buildworker-setup.service.in
[Unit]
Description=buildbot worker setup
Wants=cloud-config.target remote-fs.target local-fs.target
After=cloud-config.target remote-fs.target local-fs.target
RequiresMountsFor=@WORKERDIR@

[Service]
Type=oneshot
ExecStartPre=+@sbindir@/buildworker-init @WORKERDIR@ @WORKERUSER@
ExecStart=@sbindir@/buildworker-environment-setup @WORKERDIR@
RemainAfterExit=yes
User=@WORKERUSER@
WorkingDirectory=~
PAMName=login

[Install]
WantedBy=multi-user.target
