# src/buildworker-setup.service.in
[Unit]
Description=buildbot worker setup
Wants=cloud-config.target remote-fs.target local-fs.target
After=cloud-config.target remote-fs.target local-fs.target
RequiresMountsFor=@WORKERDIR@

[Service]
Type=oneshot
ExecStart=@sbindir@/buildworker-init @WORKERDIR@ @WORKERUSER@
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
